# openobscure-core/install/systemd/openobscure.service
[Unit]
Description=OpenObscure PII Privacy Proxy
Documentation=https://github.com/nicholasgasior/OpenObscure
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/openobscure serve
Restart=on-failure
RestartSec=5
Environment=OPENOBSCURE_CONFIG=/etc/openobscure/openobscure.toml
WorkingDirectory=/var/lib/openobscure

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/openobscure /var/log/openobscure
PrivateTmp=true

# Resource limits
MemoryMax=275M
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
